# upstart/status-led/server.conf
description "status-led server"
author "Andrew Nicolaou <andrew.nicolaou@bbc.co.uk>"

start on startup
stop on shutdown

exec /usr/local/status_led/bin/status-led-server /run/status-led
respawn
