# upstart/status-led/wifi-config.conf
description "status-led: system starting state"
author "Andrew Nicolaou <andrew.nicolaou@bbc.co.uk>"

start on wifi-config
 
task

exec /usr/local/status_led/bin/status-led /run/status-led config
